# scripts/daily-reboot.service
[Unit]
Description=Daily reboot of Raspberry Pi

[Service]
Type=oneshot
ExecStart=/sbin/reboot